# ops/roam/launchd/dev.snowline.roam.platform.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the roam platform (replication-continuity §5.1).
  Install: edit the two PATHs below, then
    cp dev.snowline.roam.platform.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/dev.snowline.roam.platform.plist
  KeepAlive restarts it after a crash; loopback bind is enforced by
  run-service.sh (--host 127.0.0.1). Repeat for governance + memory.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.snowline.roam.platform</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- EDIT: absolute path to ops/roam/run-service.sh in your checkout -->
        <string>/Users/CHANGEME/Projects/Snowline/snowline/ops/roam/run-service.sh</string>
        <string>platform</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- EDIT: absolute path to your filled-in env file -->
        <key>SNOWLINE_ENV_FILE</key>
        <string>/Users/CHANGEME/.config/snowline/env.roam</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/snowline-roam-platform.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/snowline-roam-platform.err.log</string>
</dict>
</plist>
